# dynamic-layers/stm-st-stm32mp/recipes-qt/examples/phytec-qtdemo/phytec-qtdemo.service
[Unit]
Description=Autostart Qt 5 Demo
Wants=systemd-udev-settle.service
After=systemd-udev-settle.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/qtLauncher /usr/bin/QtDemo

[Install]
WantedBy=multi-user.target
